# scripts/supervisor_image_worker.conf
[program:criavibe_image_worker]
command=php /var/www/criavibe/api/workers/image_worker.php
directory=/var/www/criavibe
autostart=true
autorestart=true
user=www-data
stdout_logfile=/var/log/supervisor/criavibe_image_worker.log
stderr_logfile=/var/log/supervisor/criavibe_image_worker.err
environment=HOME="/var/www/criavibe",USER="www-data"
